# rootfs/etc/services.d/servuo/run
#!/bin/execlineb
cd /opt/ServUO
env DOTNET_CLI_HOME=/opt/ServUO mono ServUO.exe -debug -noconsole