# rootfs/etc/services.d/servuo/run
#!/bin/execlineb
cd /opt/ServUO
mono ServUO.exe -noconsole