# rootfs/etc/services.d/servuo/run
#!/bin/execlineb
cd /opt/ServUO
export DOTNET_SYSTEM_GLOBALIZATION_INVARIANT=1
export PATH=$PATH:$HOME/.dotnet:$HOME/.dotnet/tools
mono ServUO.exe -debug -noconsole